# debian/postinst
#!/bin/sh
set -e

HLE_USER="hle-client"
HLE_HOME="/var/lib/hle-client"
VENV_DIR="/opt/hle-client/venv"

case "$1" in
    configure)
        # Create the system user and group if they do not exist yet
        if ! getent group "$HLE_USER" > /dev/null 2>&1; then
            addgroup --system "$HLE_USER"
        fi
        if ! getent passwd "$HLE_USER" > /dev/null 2>&1; then
            adduser --system \
                    --ingroup "$HLE_USER" \
                    --home "$HLE_HOME" \
                    --no-create-home \
                    --shell /usr/sbin/nologin \
                    --disabled-password \
                    --gecos "HLE Client service account" \
                    "$HLE_USER"
        fi

        # Create home directory with appropriate permissions
        if [ ! -d "$HLE_HOME" ]; then
            mkdir -p "$HLE_HOME"
        fi
        chown "$HLE_USER:$HLE_USER" "$HLE_HOME"
        chmod 750 "$HLE_HOME"

        # Create the Python virtual environment and install hle-client from PyPI
        VERSION=$(dpkg-query -W -f='${Version}' hle-client 2>/dev/null | sed 's/-.*//')
        if [ ! -d "$VENV_DIR" ] || [ ! -x "$VENV_DIR/bin/hle" ]; then
            python3 -m venv "$VENV_DIR"
            "$VENV_DIR/bin/pip" install --quiet "hle-client==$VERSION"
        elif [ -x "$VENV_DIR/bin/pip" ]; then
            # Upgrade existing installation to the new version
            "$VENV_DIR/bin/pip" install --quiet --upgrade "hle-client==$VERSION"
        fi
        chown -R "$HLE_USER:$HLE_USER" "$VENV_DIR"

        # Create config directory and a template config file under the user's home
        CONFIG_DIR="$HLE_HOME/.config/hle"
        CONFIG_FILE="$CONFIG_DIR/config.toml"
        mkdir -p "$CONFIG_DIR"
        if [ ! -f "$CONFIG_FILE" ]; then
            cat > "$CONFIG_FILE" << 'TOML'
# HLE Client configuration
#
# Obtain your API key from https://app.hle.world and set it here,
# or run the following command to authenticate interactively:
#   sudo -u hle-client /opt/hle-client/venv/bin/hle auth login
#
# api_key = "your-api-key-here"
TOML
            chmod 600 "$CONFIG_FILE"
        fi
        # Apply ownership and permissions after all file operations
        chown -R "$HLE_USER:$HLE_USER" "$HLE_HOME/.config"
        chmod 700 "$CONFIG_DIR"

        # Reload systemd unit files
        if [ -d /run/systemd/system ]; then
            systemctl daemon-reload || true
        fi

        echo ""
        echo "HLE Client installed successfully."
        echo ""
        echo "Next steps:"
        echo "  1. Set the service URL in /etc/default/hle-client:"
        echo "       sudo editor /etc/default/hle-client"
        echo "     Uncomment and set: HLE_SERVICE_URL=http://localhost:8080"
        echo ""
        echo "  2. Authenticate (or set HLE_API_KEY in /etc/default/hle-client):"
        echo "       sudo -u hle-client /opt/hle-client/venv/bin/hle auth login"
        echo ""
        echo "  3. Enable and start the service:"
        echo "       sudo systemctl enable --now hle-client"
        echo ""
        ;;

    abort-upgrade|abort-remove|abort-deconfigure)
        ;;

    *)
        echo "postinst called with unknown argument '$1'" >&2
        exit 1
        ;;
esac

#DEBHELPER#

exit 0
